# roles/systemd_rc_local/templates/rc.local.j2
#!/bin/bash
#
# {{ ansible_managed }}

touch /var/lock/subsys/local

{% for cmd in systemd_rc_local_list %}
{{ cmd }}
{% endfor %}
