# roles/rc_local/templates/rc.local.j2
#!/bin/bash
#
# {{ ansible_managed }}

touch /var/lock/subsys/local

{% for cmd in rc_local_list %}
{{ cmd }}
{% endfor %}
